# config/supervisor/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"

[eventlistener:letsencrypt]
command=/opt/letsencrypt-run.py
events=TICK_3600
